# services/chord-bot-api.service
[Unit]
Description=Chord-Bot API (FastAPI/Flask)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{RUN_AS_USER}}
WorkingDirectory={{CHORD_BOT_API_DIR}}
ExecStart={{CHORD_BOT_API_VENV}}/bin/bash -lc '{{CHORD_BOT_API_START_CMD}}'
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1
StandardOutput=append:{{CHORD_BOT_API_LOG_DIR}}/stdout.log
StandardError=append:{{CHORD_BOT_API_LOG_DIR}}/stderr.log

[Install]
WantedBy=multi-user.target
