# services/chord-bot.service
[Unit]
Description=Chord-Bot local app
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{RUN_AS_USER}}
WorkingDirectory={{CHORD_BOT_DIR}}
ExecStart={{CHORD_BOT_VENV}}/bin/bash -lc '{{CHORD_BOT_START_CMD}}'
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1
StandardOutput=append:{{CHORD_BOT_LOG_DIR}}/stdout.log
StandardError=append:{{CHORD_BOT_LOG_DIR}}/stderr.log

[Install]
WantedBy=multi-user.target
